# ansible/roles/hashi_vault/templates/etc/profile.d/vault.sh
#!/bin/sh
export VAULT_ADDR=http://192.168.77.25:8200
